# usr/lib/systemd/system/cachyos-iw-set-regdomain.path
# SPDX-License-Identifier: ISC
# SPDX-FileCopyrightText: 2025 Thomas Duckworth <tduck@filotimoproject.org>

[Unit]
Description=Monitor Timezone Changes to Set Wireless Regulatory Domain

[Path]
PathChanged=/etc/localtime
Unit=cachyos-iw-set-regdomain.service

[Install]
WantedBy=multi-user.target
