# etc/udev/rules.d/10-laptop.rules
ACTION=="add", SUBSYSTEM=="power_supply", KERNEL!="hidpp_battery*", RUN+="/usr/bin/laptop-mode"
